# Linux/etc/profile.d/custom-aliases.sh
#===============================================================================
# SCRIPT HEADER TEMPLATE
#============================
# Script Name: custom-aliases.sh
# Description: Creates shell aliases for the all users
# 
# Created By: robinhood1995@yahoo.com
# Created On: 2025-12-18
# 
# Modified By: 
# Modified On: 
# Modification Notes:
# 
# Version: 1.0.0
# License: MIT License
# 
# Dependencies: 
# Usage: Auto configures the shell when logging in
#
# WARNING:
# 
# Changelog:
# - [YYYY-MM-DD]: [VERSION] - [SUMMARY OF CHANGES]
# - [YYYY-MM-DD]: [VERSION] - [INITIAL CREATION]
#===============================================================================

# Common interactive aliases for all users

alias rm='rm -i'
alias cp='cp -i'
alias mv='mv -i'
alias vi='vim'
alias tailf='tail -f'